# backend/systemd/talky-api.service
[Unit]
Description=Talky.ai API Server
After=network.target redis.service
Requires=redis.service

[Service]
Type=exec
WorkingDirectory=/home/ai-lab/Desktop/Talky.ai-complete-/backend
EnvironmentFile=/home/ai-lab/Desktop/Talky.ai-complete-/backend/.env
ExecStart=/home/ai-lab/Desktop/Talky.ai-complete-/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4 --ws websockets --ws-ping-interval 30 --ws-ping-timeout 5
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=talky-api

[Install]
WantedBy=multi-user.target
