# backend/systemd/talky-reminder-worker.service
[Unit]
Description=Talky.ai Reminder Worker
After=network.target talky-api.service

[Service]
Type=exec
WorkingDirectory=/opt/talky/backend
EnvironmentFile=/opt/talky/backend/.env
ExecStart=/opt/talky/backend/venv/bin/python -m app.workers.reminder_worker
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=talky-reminder

[Install]
WantedBy=multi-user.target
