# backend/systemd/talky-dialer-worker.service
[Unit]
Description=Talky.ai Dialer Worker
After=network.target talky-api.service

[Service]
Type=exec
WorkingDirectory=/opt/talky/backend
EnvironmentFile=/opt/talky/backend/.env
ExecStart=/opt/talky/backend/venv/bin/python -m app.workers.dialer_worker
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=talky-dialer

[Install]
WantedBy=multi-user.target
