# backend/systemd/talky-api.service
[Unit]
Description=Talky.ai API Server
After=network.target

[Service]
Type=exec
WorkingDirectory=/opt/talky/backend
EnvironmentFile=/opt/talky/backend/.env
ExecStart=/opt/talky/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 2 --ws websockets --ws-ping-interval 30 --ws-ping-timeout 5
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=talky-api

[Install]
WantedBy=multi-user.target
