# backend/systemd/talky-reminder-worker.service
[Unit]
Description=Talky.ai Reminder Worker
After=network.target talky-api.service
# Crash-loop guard: cap restarts at 5 per 300s so a persistently broken deploy
# stops (enters failed state) instead of thrashing forever. These are [Unit]
# directives in modern systemd (>=230); placing them under [Service] is ignored.
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
# Type=notify: systemd waits for the worker's READY=1 (sent from the heartbeat
# loop at startup) before considering the unit started. WatchdogSec=180 arms
# the watchdog — the worker pets it with WATCHDOG=1 every ~60s from the same
# event loop that runs the reminder-scan loop, so a HUNG (not just dead)
# worker whose loop stops petting is force-restarted. NotifyAccess=main
# restricts the notify socket to the main process.
Type=notify
NotifyAccess=main
WatchdogSec=180
WorkingDirectory=/opt/talky/backend
EnvironmentFile=/opt/talky/backend/.env
ExecStart=/opt/talky/backend/venv/bin/python -m app.workers.reminder_worker
# Restart=always so a clean sys.exit(1) (fatal consecutive-error path) or a
# watchdog kill both bring the worker back (crash-loop capped via [Unit]
# StartLimit* above).
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=talky-reminder

[Install]
WantedBy=multi-user.target
