# backend/deploy/systemd/talky-trunk-status.service
[Unit]
Description=Talky SIP-trunk live registration status updater (reads pjsip show registrations -> DB)
After=asterisk.service postgresql.service

[Service]
Type=oneshot
# root: needs the asterisk CLI + reads /opt/talky/backend/.env for DATABASE_URL
User=root
ExecStart=/opt/talky/backend/venv/bin/python /opt/talky/backend/scripts/trunk_live_status_updater.py
TimeoutStartSec=30
